# src/profile.d/bash-insulter.sh
[ -f /etc/bash.command-not-found-messages ] && . /etc/bash.command-not-found-messages

[ -f /etc/bash.command-not-found ] && . /etc/bash.command-not-found